# squashfs-root/etc/init.d/safe.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

SERVICE_USE_PID=1

START=50
kver=`ls /lib/modules`
start() {
	insmod /lib/modules/${kver}/safe.ko
}

stop() {
	rmmod /lib/modules/${kver}/safe.ko
}


