# squashfs-root/etc/init.d/sqos
#!/bin/sh /etc/rc.common

SERVICE_USE_PID=1
START=80

kver=`ls /lib/modules`
start() {
	insmod /lib/modules/${kver}/sqos.ko
}

